# board/systemd/deepx-dual.service
[Unit]
Description=DeepX Dual Demo (SCDepthV3 + YOLO26m-cls on DeepX, TrOCR on HTP)
Documentation=https://github.com/dvirchakim/QCM8550_8550-deepx_m1m
Wants=init_qti_graphics.service init_display.service qmmf-server.service
After=init_qti_graphics.service init_display.service qmmf-server.service

[Service]
Type=simple
User=root
Environment=ADSP_LIBRARY_PATH=/system/lib/rfsa/adsp;/system/vendor/lib/rfsa/adsp;/dsp
Environment=XDG_RUNTIME_DIR=/run/user/root
Environment=WAYLAND_DISPLAY=wayland-1
Environment=QT_QPA_PLATFORM=wayland-egl
Environment=QT_WAYLAND_SHELL_INTEGRATION=wl-shell

ExecStartPre=/bin/sh -c 'systemctl stop yolo26-parallel.service imdt-deepx-demo.service 2>/dev/null || true; pkill -9 dual_deepx_worker trocr_worker seg_worker yolo26_launcher gst-launch-1.0 2>/dev/null || true; sleep 2'
ExecStartPre=/bin/sh -c 'cd /lib/modules/$(uname -r)/extra && insmod dxrt_driver.ko 2>/dev/null || true'
ExecStart=/bin/sh /data/local/tmp/deepx_dual.sh
ExecStop=/bin/sh -c 'pkill -9 dual_deepx_worker trocr_worker gst-launch-1.0 2>/dev/null || true'

Restart=on-failure
RestartSec=8
TimeoutStartSec=120
TimeoutStopSec=5
KillMode=control-group

StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepx-dual
